# configs/pbc.conf
[program:pbc]
command=python3.6 main.py
directory = /root/postback_collector/app
autostart = true
stderr_logfile = /var/log/postback_collector.log
stdout_logfile = /var/log/postback_collector.log