# configs/pbc_redis.conf
[program:pbc_redis]
command=redis-server
autostart = true
stderr_logfile = /var/log/postback_collector_redis.log
stdout_logfile = /var/log/postback_collector_redis.log